# wttrTimer.timer
[Unit]
Description=Daily Weather Report Timer

[Timer]
Unit=wttrScript.service
OnBootSec=15min
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
